# workers/agent/supervisord/foyer-agent-worker.conf
; FastAPI agent worker for Foyer. Two processes from day one — agent
; turns are LLM-bound and benefit from horizontal pickup. Supervisord
; restarts on crash; stderr + stdout funnel into per-process logs under
; /var/log/foyer/.
;
; Deploys to /etc/supervisor/conf.d/foyer-agent-worker.conf on the EC2 box.
; Reload with:
;   sudo supervisorctl reread && sudo supervisorctl update

[program:foyer-agent-worker]
command=/var/www/foyer/current/workers/agent/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 88%(process_num)02d --workers 1
directory=/var/www/foyer/current/workers/agent
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
user=ubuntu
environment=FOYER_INTERNAL_SECRET="REPLACE_ME",FOYER_LARAVEL_INTERNAL_BASE_URL="http://127.0.0.1:8000",FOYER_LLM_PROVIDER="mock",FOYER_LLM_MODEL="claude-sonnet-4-6",FOYER_DEFAULT_COST_CEILING_MICROS="500000"
stdout_logfile=/var/log/foyer/agent-%(process_num)02d.out.log
stderr_logfile=/var/log/foyer/agent-%(process_num)02d.err.log
stopasgroup=true
killasgroup=true
priority=20
