# infra/supervisor/foyer-horizon.conf
; Foyer Horizon worker. Deploys to /etc/supervisor/conf.d/foyer-horizon.conf
; on the EC2 host. Reload with:
;   sudo supervisorctl reread && sudo supervisorctl update

[program:foyer-horizon]
process_name=%(program_name)s
command=php /var/www/foyer/current/artisan horizon
autostart=true
autorestart=true
user=ubuntu
redirect_stderr=true
stdout_logfile=/var/log/foyer/horizon.log
stopwaitsecs=3600
